# platform/broadcom/sonic-platform-modules-dell/s3000/systemd/platform-modules-s3000.service
[Unit]
Description=Dell S3000 Platform modules
Before=pmon.service
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/local/bin/s3000_platform.sh init
ExecStop=/usr/local/bin/s3000_platform.sh deinit
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
