# files/image_config/interfaces/interfaces-config.service
[Unit]
Description=Update interfaces configuration
Requires=updategraph.service
After=updategraph.service
ConditionPathExists=!/usr/bin/pde.sh

[Service]
Type=oneshot
ExecStart=/usr/bin/interfaces-config.sh

[Install]
WantedBy=multi-user.target
